# roles/cryptshot/templates/cryptshot.service.j2
Description=Cryptshot Backup Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cryptshot-execute
StandardOutput=syslog
StandardError=syslog
