# roles/cryptshot/templates/cryptshot.timer.j2
[Unit]
Description=Cryptshot Backup Timer

[Timer]
Unit=cryptshot.service
OnCalendar=hourly
Persistent=True

[Install]
WantedBy=timers.target
