# roles/cryptshot/templates/cryptshot.service.j2
Description=Cryptshot Backup Service: {{ target.name|capitalize }}

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cryptshot-{{ target.name }}
StandardOutput=syslog
StandardError=syslog
