# roles/himawaripy/templates/himawaripy.service.j2
# {{ ansible_managed }}
[Unit]
Description=Update desktop background with satellite imagery

[Service]
Type=oneshot
ExecStart=/usr/local/bin/himawaripy {{ himawaripy.flags }}
