# roles/laptop/files/chargecontrol.service
[Service]
ExecStart=ectool chargecontrol normal 87 87

[Install]
WantedBy=default.target
