# roles/laptop/files/chargecontrol.service
[Service]
ExecStart=framework_tool --charge-limit 87

[Install]
WantedBy=default.target
