# roles/laptop/templates/lowbatt.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Low battery warning timer

[Timer]
Unit=lowbatt.service
OnCalendar=*:0/2

[Install]
WantedBy=timers.target
