# roles/base/files/paccache.service
[Unit]
Description=Clean-up Old Pacman Files

[Service]
Type=oneshot
IOSchedulingClass=idle
CPUSchedulingPolicy=idle
ExecStart=/usr/bin/paccache -r
ExecStart=/usr/bin/paccache -ruk0
