# roles/base/files/paccache.timer
[Unit]
Description=Clean-up Old Pacman Files

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3h

[Install]
WantedBy=multi-user.target
