# roles/mirrorlist/files/reflector-update.service
[Unit]
Description=Reflector Update Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/reflector-update
