# roles/backup/templates/tarsnapper.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Tarsnapper Synchronization Timer

[Timer]
OnUnitActiveSec={{ tarsnapper.timer.frequency }}
Unit=tarsnapper.service

[Install]
WantedBy=timers.target
