# roles/units/templates/units_cur.timer.j2
# {{ ansible_managed }}
[Unit]
Description=GNU Units Currency Update Timer

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3600
Unit=units_cur.service

[Install]
WantedBy=timers.target
